# launchd/com.idp.xero-cash-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.idp.xero-cash-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/yostenterprises/IDP_Command_Center/agents/xero_cash_agent.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/yostenterprises/IDP_Command_Center/logs/xero_cash.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yostenterprises/IDP_Command_Center/logs/xero_cash_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
